# cfg/rc.local
#
# ХЕРНЯ!
#
# надо прописывать полные пути до исполняемых файлов, иначе не работает
#

# start gitit
cd ~/XFiles/gitit && eval "$(which gitit) &"
